# roles/rclone/templates/rclone-job@.timer.j2
[Timer]
OnCalendar=
OnCalendar={{ item.calendar }}
