# roles/autofs/templates/service.automount.j2
[Unit]
Description=AutoFS {{ item.name }}

[Automount]
Where={{ item.mountpoint }}

[Install]
WantedBy=multi-user.target
