# roles/nas/files/zfs-scrub@.service
[Unit]
Description=Scrub ZFS tank pool

[Service]
Type=oneshot
ExecStartPre=-/usr/bin/zpool scrub -s %i
ExecStart=/usr/bin/zpool scrub %i
